# services/reset-backlight.service
[Unit]
Description=Dim screen on resume.
After=suspend.target

[Service]
ExecStart=/usr/bin/xbacklight -set 4

[Install]
WantedBy=suspend.target
